# infra/systemd/auto_update.timer
[Unit]
Description=Daily timer for SnagDeals auto-update

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
